# agent/includes/ogp_agent.init.dbn
#!/bin/bash
#
### BEGIN INIT INFO
# Provides:          ogp_agent
# Required-Start:    $network $local_fs
# Required-Stop:     $local_fs $network
# Should-Start:      $local_fs $network 
# Should-Stop:       $local_fs $network 
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Short-Description: Start and stop the OGP Agent
# Description:       Start and stop the OGP Agent
### END INIT INFO
#
set -e
set -u
${DEBIAN_SCRIPT_DEBUG:+ set -v -x}

. /lib/lsb/init-functions

agent_dir=OGP_AGENT_DIR
agent_user=OGP_USER

#
# main()
#

if [ "$( whoami )" != "root" ]
then
	if [ -f "/usr/bin/sudo" ] && [ "$( groups $agent_user | grep "\bsudo\b" )" != "" ]
	then
		sudo /etc/init.d/ogp_agent ${1:-''}
		exit
	else
		echo "Permission denied."
		exit
	fi
fi

case "${1:-''}" in
    'start')
	# Lets the agent user to use sudo to enable FTP accounts and use renice and taskset.
	if [ "$( groups $agent_user | grep "\bsudo\b" )" == "" ]
	then
		adduser $agent_user sudo &> /dev/null
	fi
	
	group=`groups $agent_user | awk '{ print $3 }'`; 
	chown -Rf $agent_user:$group $agent_dir &> /dev/null

	# Lets the agent user to attach screens.
	if [ "$( groups $agent_user | grep "\btty\b" )" == "" ]
	then
		adduser $agent_user tty &> /dev/null
	fi
	chmod g+rw /dev/pts/* &> /dev/null
	chmod g+rw /dev/tty* &> /dev/null
	
	# Check the FTP status
	if [ -f "/etc/init.d/pure-ftpd" ]
    then
		if [ "$( service pure-ftpd status | grep "* pure-ftpd is not running" )" == "" ]
		then
			service pure-ftpd stop &> /dev/null
			if [ "$( service pure-ftpd start | grep "pureftpd.pdb" )" == "" ]
			then
				service pure-ftpd stop &> /dev/null
				echo no > /etc/pure-ftpd/conf/PAMAuthentication
				echo no > /etc/pure-ftpd/conf/UnixAuthentication
				echo yes > /etc/pure-ftpd/conf/CreateHomeDir
				touch /etc/pure-ftpd/pureftpd.passwd
				ln -s /etc/pure-ftpd/pureftpd.passwd /etc/pureftpd.passwd
				ln -s /etc/pure-ftpd/conf/PureDB /etc/pure-ftpd/auth/50pure
				ln -s /etc/pure-ftpd/pureftpd.pdb /etc/pureftpd.pdb
				pure-pw mkdb &> /dev/null
				service pure-ftpd start &> /dev/null
			fi
		elif [ "$( service pure-ftpd start | grep "pureftpd.pdb" )" == "" ]
		then
			service pure-ftpd stop &> /dev/null
			echo no > /etc/pure-ftpd/conf/PAMAuthentication
			echo no > /etc/pure-ftpd/conf/UnixAuthentication
			echo yes > /etc/pure-ftpd/conf/CreateHomeDir
			touch /etc/pure-ftpd/pureftpd.passwd
			ln -s /etc/pure-ftpd/pureftpd.passwd /etc/pureftpd.passwd
			ln -s /etc/pure-ftpd/conf/PureDB /etc/pure-ftpd/auth/50pure
			ln -s /etc/pure-ftpd/pureftpd.pdb /etc/pureftpd.pdb
			pure-pw mkdb &> /dev/null
			service pure-ftpd start &> /dev/null		
		fi
	fi
	
    cd $agent_dir
    su -c "screen -d -m -t ogp_agent -c ogp_screenrc -S ogp_agent ./ogp_agent_run -pidfile ogp_agent_run.pid" $agent_user &> $agent_dir/ogp_agent.svc &
    log_daemon_msg "Starting OGP Agent" "ogp_agent_run -pidfile ogp_agent_run.pid"
    ;;

    'stop')
    log_daemon_msg "Stopping OGP Agent" "ogp_agent_run"
    kill `cat $agent_dir/ogp_agent_run.pid`
    ;;

    'restart')
    set +e; /etc/init.d/ogp_agent stop; set -e
    /etc/init.d/ogp_agent start
    ;;

    *)
    echo "Usage: /etc/init.d/ogp_agent start|stop|restart"
    exit 1
    ;;
esac

exit 0;
